# supervisord.conf
[supervisord]
nodaemon=true

[program:pushion]
command=/sbin/my_init
autorestart=true

[program:plirocron]
command=cron
